# tuned.upstart.conf
# tuned - tuned
#

description "tuned"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/sbin/tuned
